# extras/etc/udev/rules.d/99-battery.rules
SUBSYSTEM=="power_supply", RUN+="/usr/bin/su john -c /usr/local/bin/notify-battery-status.sh"
